# roles/bootstrap_docker/templates/etc/systemd/system/docker.service.d/environment.conf.j2
{{ 'Ansible managed' | comment(format="# {}") }}

[Service]
Environment="{{ bootstrap_docker__daemon_environment | join('" "') }}"
