# roles/bootstrap_docker/templates/etc/systemd/system/docker.service.d/options.conf.j2
{{ 'Ansible managed' | comment(format="# {}") }}

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd {{ bootstrap_docker__daemon_flags | join(" ") }}
