# roles/bootstrap_llm_host/templates/nemoclaw.service.j2
[Unit]
Description=NemoClaw Autonomous Agent Orchestration Daemon
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User={{ bootstrap_llm_host__nemoclaw_user }}
Group={{ bootstrap_llm_host__nemoclaw_user }}
WorkingDirectory={{ bootstrap_llm_host__nemoclaw_home_dir }}
ExecStart=/usr/local/bin/nemoclaw serve --config {{ bootstrap_llm_host__nemoclaw_home_dir }}/.config/nemoclaw/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
