# roles/bootstrap-docker/templates/etc/systemd/system/docker.service.d/environment.conf.j2
# {{ ansible_managed }}

[Service]
Environment="{{ bootstrap_docker__daemon_environment | join('" "') }}"
